# docker/supervisord.conf
[supervisord]
nodaemon = true
logfile = /dev/stdout
logfile_maxbytes = 0

[program:php-fpm]
command = /usr/local/sbin/php-fpm -F

[program:nginx]
command = /usr/sbin/nginx -g 'daemon off;'